# roles/btsync/files/btsync.conf
# btsync
#

description	"BTSync"

start on startup
stop on runlevel [!023456]

#console output

exec sudo -u ubuntu /opt/btsync/btsync --nodaemon 2>&1 | logger -t "btsync"

