# roles/btsync/files/btsync.conf
# btsync
#

description	"BTSync"

start on startup
stop on runlevel [!023456]

#console output

exec /opt/btsync/btsync --nodaemon | logger -t "BTSync"

